# guardian/config/supervisor/sy-ai.conf
[program:ai_OTG]
autorestart = true
redirect_stderr = true
stopsignal = QUIT
pidfile = /var/safely-you/run/guardian/ai_OTG.pid
command = /opt/safely-you/guardian/venv/bin/python manage.py services ai run --ai-type=OTG
user = safelyyou
autostart = true
directory = /opt/safely-you/guardian
stdout_logfile = /var/safely-you/logs/guardian/ai.log
exitcodes = 0
