# guardian/config/supervisor/sy-streams.conf
[program:stream-00.13.E2.27.0F.A9]
autorestart = true
redirect_stderr = true
stopsignal = QUIT
pidfile = /var/safely-you/run/guardian/stream-00.13.E2.27.0F.A9.pid
command = /opt/safely-you/guardian/venv/bin/python manage.py services stream 00:13:E2:27:0F:A9 run
user = safelyyou
autostart = true
directory = /opt/safely-you/guardian
stdout_logfile = /var/safely-you/logs/guardian/stream-00.13.E2.27.0F.A9.log
exitcodes = 0

[program:stream-00.13.E2.27.0F.B5]
autorestart = true
redirect_stderr = true
stopsignal = QUIT
pidfile = /var/safely-you/run/guardian/stream-00.13.E2.27.0F.B5.pid
command = /opt/safely-you/guardian/venv/bin/python manage.py services stream 00:13:E2:27:0F:B5 run
user = safelyyou
autostart = true
directory = /opt/safely-you/guardian
stdout_logfile = /var/safely-you/logs/guardian/stream-00.13.E2.27.0F.B5.log
exitcodes = 0

[program:stream-00.13.E2.27.0F.F1]
autorestart = true
redirect_stderr = true
stopsignal = QUIT
pidfile = /var/safely-you/run/guardian/stream-00.13.E2.27.0F.F1.pid
command = /opt/safely-you/guardian/venv/bin/python manage.py services stream 00:13:E2:27:0F:F1 run
user = safelyyou
autostart = true
directory = /opt/safely-you/guardian
stdout_logfile = /var/safely-you/logs/guardian/stream-00.13.E2.27.0F.F1.log
exitcodes = 0

[program:stream-00.13.E2.27.10.51]
autorestart = true
redirect_stderr = true
stopsignal = QUIT
pidfile = /var/safely-you/run/guardian/stream-00.13.E2.27.10.51.pid
command = /opt/safely-you/guardian/venv/bin/python manage.py services stream 00:13:E2:27:10:51 run
user = safelyyou
autostart = true
directory = /opt/safely-you/guardian
stdout_logfile = /var/safely-you/logs/guardian/stream-00.13.E2.27.10.51.log
exitcodes = 0

[program:stream-00.13.E2.27.14.F3]
autorestart = true
redirect_stderr = true
stopsignal = QUIT
pidfile = /var/safely-you/run/guardian/stream-00.13.E2.27.14.F3.pid
command = /opt/safely-you/guardian/venv/bin/python manage.py services stream 00:13:E2:27:14:F3 run
user = safelyyou
autostart = true
directory = /opt/safely-you/guardian
stdout_logfile = /var/safely-you/logs/guardian/stream-00.13.E2.27.14.F3.log
exitcodes = 0

[program:stream-00.13.E2.27.15.2D]
autorestart = true
redirect_stderr = true
stopsignal = QUIT
pidfile = /var/safely-you/run/guardian/stream-00.13.E2.27.15.2D.pid
command = /opt/safely-you/guardian/venv/bin/python manage.py services stream 00:13:E2:27:15:2D run
user = safelyyou
autostart = true
directory = /opt/safely-you/guardian
stdout_logfile = /var/safely-you/logs/guardian/stream-00.13.E2.27.15.2D.log
exitcodes = 0

[program:stream-00.13.E2.27.15.4B]
autorestart = true
redirect_stderr = true
stopsignal = QUIT
pidfile = /var/safely-you/run/guardian/stream-00.13.E2.27.15.4B.pid
command = /opt/safely-you/guardian/venv/bin/python manage.py services stream 00:13:E2:27:15:4B run
user = safelyyou
autostart = true
directory = /opt/safely-you/guardian
stdout_logfile = /var/safely-you/logs/guardian/stream-00.13.E2.27.15.4B.log
exitcodes = 0

[program:stream-00.13.E2.27.15.53]
autorestart = true
redirect_stderr = true
stopsignal = QUIT
pidfile = /var/safely-you/run/guardian/stream-00.13.E2.27.15.53.pid
command = /opt/safely-you/guardian/venv/bin/python manage.py services stream 00:13:E2:27:15:53 run
user = safelyyou
autostart = true
directory = /opt/safely-you/guardian
stdout_logfile = /var/safely-you/logs/guardian/stream-00.13.E2.27.15.53.log
exitcodes = 0

[program:stream-00.13.E2.27.16.86]
autorestart = true
redirect_stderr = true
stopsignal = QUIT
pidfile = /var/safely-you/run/guardian/stream-00.13.E2.27.16.86.pid
command = /opt/safely-you/guardian/venv/bin/python manage.py services stream 00:13:E2:27:16:86 run
user = safelyyou
autostart = true
directory = /opt/safely-you/guardian
stdout_logfile = /var/safely-you/logs/guardian/stream-00.13.E2.27.16.86.log
exitcodes = 0

[program:stream-00.13.E2.27.16.9D]
autorestart = true
redirect_stderr = true
stopsignal = QUIT
pidfile = /var/safely-you/run/guardian/stream-00.13.E2.27.16.9D.pid
command = /opt/safely-you/guardian/venv/bin/python manage.py services stream 00:13:E2:27:16:9D run
user = safelyyou
autostart = true
directory = /opt/safely-you/guardian
stdout_logfile = /var/safely-you/logs/guardian/stream-00.13.E2.27.16.9D.log
exitcodes = 0

[program:stream-00.13.E2.27.17.D4]
autorestart = true
redirect_stderr = true
stopsignal = QUIT
pidfile = /var/safely-you/run/guardian/stream-00.13.E2.27.17.D4.pid
command = /opt/safely-you/guardian/venv/bin/python manage.py services stream 00:13:E2:27:17:D4 run
user = safelyyou
autostart = true
directory = /opt/safely-you/guardian
stdout_logfile = /var/safely-you/logs/guardian/stream-00.13.E2.27.17.D4.log
exitcodes = 0

[program:stream-00.13.E2.27.17.D6]
autorestart = true
redirect_stderr = true
stopsignal = QUIT
pidfile = /var/safely-you/run/guardian/stream-00.13.E2.27.17.D6.pid
command = /opt/safely-you/guardian/venv/bin/python manage.py services stream 00:13:E2:27:17:D6 run
user = safelyyou
autostart = true
directory = /opt/safely-you/guardian
stdout_logfile = /var/safely-you/logs/guardian/stream-00.13.E2.27.17.D6.log
exitcodes = 0

[program:stream-00.13.E2.27.17.EC]
autorestart = true
redirect_stderr = true
stopsignal = QUIT
pidfile = /var/safely-you/run/guardian/stream-00.13.E2.27.17.EC.pid
command = /opt/safely-you/guardian/venv/bin/python manage.py services stream 00:13:E2:27:17:EC run
user = safelyyou
autostart = true
directory = /opt/safely-you/guardian
stdout_logfile = /var/safely-you/logs/guardian/stream-00.13.E2.27.17.EC.log
exitcodes = 0

[program:stream-00.13.E2.27.17.FF]
autorestart = true
redirect_stderr = true
stopsignal = QUIT
pidfile = /var/safely-you/run/guardian/stream-00.13.E2.27.17.FF.pid
command = /opt/safely-you/guardian/venv/bin/python manage.py services stream 00:13:E2:27:17:FF run
user = safelyyou
autostart = true
directory = /opt/safely-you/guardian
stdout_logfile = /var/safely-you/logs/guardian/stream-00.13.E2.27.17.FF.log
exitcodes = 0

[program:stream-00.13.E2.27.18.10]
autorestart = true
redirect_stderr = true
stopsignal = QUIT
pidfile = /var/safely-you/run/guardian/stream-00.13.E2.27.18.10.pid
command = /opt/safely-you/guardian/venv/bin/python manage.py services stream 00:13:E2:27:18:10 run
user = safelyyou
autostart = true
directory = /opt/safely-you/guardian
stdout_logfile = /var/safely-you/logs/guardian/stream-00.13.E2.27.18.10.log
exitcodes = 0

[program:stream-00.13.E2.27.18.14]
autorestart = true
redirect_stderr = true
stopsignal = QUIT
pidfile = /var/safely-you/run/guardian/stream-00.13.E2.27.18.14.pid
command = /opt/safely-you/guardian/venv/bin/python manage.py services stream 00:13:E2:27:18:14 run
user = safelyyou
autostart = true
directory = /opt/safely-you/guardian
stdout_logfile = /var/safely-you/logs/guardian/stream-00.13.E2.27.18.14.log
exitcodes = 0

[program:stream-00.13.E2.27.18.1F]
autorestart = true
redirect_stderr = true
stopsignal = QUIT
pidfile = /var/safely-you/run/guardian/stream-00.13.E2.27.18.1F.pid
command = /opt/safely-you/guardian/venv/bin/python manage.py services stream 00:13:E2:27:18:1F run
user = safelyyou
autostart = true
directory = /opt/safely-you/guardian
stdout_logfile = /var/safely-you/logs/guardian/stream-00.13.E2.27.18.1F.log
exitcodes = 0

[program:stream-00.13.E2.27.18.23]
autorestart = true
redirect_stderr = true
stopsignal = QUIT
pidfile = /var/safely-you/run/guardian/stream-00.13.E2.27.18.23.pid
command = /opt/safely-you/guardian/venv/bin/python manage.py services stream 00:13:E2:27:18:23 run
user = safelyyou
autostart = true
directory = /opt/safely-you/guardian
stdout_logfile = /var/safely-you/logs/guardian/stream-00.13.E2.27.18.23.log
exitcodes = 0

[program:stream-00.13.E2.27.18.2F]
autorestart = true
redirect_stderr = true
stopsignal = QUIT
pidfile = /var/safely-you/run/guardian/stream-00.13.E2.27.18.2F.pid
command = /opt/safely-you/guardian/venv/bin/python manage.py services stream 00:13:E2:27:18:2F run
user = safelyyou
autostart = true
directory = /opt/safely-you/guardian
stdout_logfile = /var/safely-you/logs/guardian/stream-00.13.E2.27.18.2F.log
exitcodes = 0

[program:stream-00.13.E2.27.18.35]
autorestart = true
redirect_stderr = true
stopsignal = QUIT
pidfile = /var/safely-you/run/guardian/stream-00.13.E2.27.18.35.pid
command = /opt/safely-you/guardian/venv/bin/python manage.py services stream 00:13:E2:27:18:35 run
user = safelyyou
autostart = true
directory = /opt/safely-you/guardian
stdout_logfile = /var/safely-you/logs/guardian/stream-00.13.E2.27.18.35.log
exitcodes = 0

[program:stream-00.13.E2.27.18.DA]
autorestart = true
redirect_stderr = true
stopsignal = QUIT
pidfile = /var/safely-you/run/guardian/stream-00.13.E2.27.18.DA.pid
command = /opt/safely-you/guardian/venv/bin/python manage.py services stream 00:13:E2:27:18:DA run
user = safelyyou
autostart = true
directory = /opt/safely-you/guardian
stdout_logfile = /var/safely-you/logs/guardian/stream-00.13.E2.27.18.DA.log
exitcodes = 0

[program:stream-00.13.E2.27.18.F9]
autorestart = true
redirect_stderr = true
stopsignal = QUIT
pidfile = /var/safely-you/run/guardian/stream-00.13.E2.27.18.F9.pid
command = /opt/safely-you/guardian/venv/bin/python manage.py services stream 00:13:E2:27:18:F9 run
user = safelyyou
autostart = true
directory = /opt/safely-you/guardian
stdout_logfile = /var/safely-you/logs/guardian/stream-00.13.E2.27.18.F9.log
exitcodes = 0

[program:stream-00.13.E2.27.1A.1E]
autorestart = true
redirect_stderr = true
stopsignal = QUIT
pidfile = /var/safely-you/run/guardian/stream-00.13.E2.27.1A.1E.pid
command = /opt/safely-you/guardian/venv/bin/python manage.py services stream 00:13:E2:27:1A:1E run
user = safelyyou
autostart = true
directory = /opt/safely-you/guardian
stdout_logfile = /var/safely-you/logs/guardian/stream-00.13.E2.27.1A.1E.log
exitcodes = 0

[program:stream-60.23.A4.BD.B9.E1]
autorestart = true
redirect_stderr = true
stopsignal = QUIT
pidfile = /var/safely-you/run/guardian/stream-60.23.A4.BD.B9.E1.pid
command = /opt/safely-you/guardian/venv/bin/python manage.py services stream 60:23:A4:BD:B9:E1 run
user = safelyyou
autostart = true
directory = /opt/safely-you/guardian
stdout_logfile = /var/safely-you/logs/guardian/stream-60.23.A4.BD.B9.E1.log
exitcodes = 0

[program:stream-50.13.95.A2.78.72]
autorestart = true
redirect_stderr = true
stopsignal = QUIT
pidfile = /var/safely-you/run/guardian/stream-50.13.95.A2.78.72.pid
command = /opt/safely-you/guardian/venv/bin/python manage.py services stream 50:13:95:A2:78:72 run
user = safelyyou
autostart = true
directory = /opt/safely-you/guardian
stdout_logfile = /var/safely-you/logs/guardian/stream-50.13.95.A2.78.72.log
exitcodes = 0
